# arch/xu4/extra/etc/systemd/system/waggle-heartbeat.service
[Unit]
Description=Waggle Heartbeat

[Service]
Type=oneshot
ExecStart=/root/waggle-heartbeat

[Install]
WantedBy=multi-user.target
