# arch/c1/extra/etc/systemd/system/waggle-heartbeat.service
[Unit]
Description=Waggle Heartbeat

[Service]
ExecStart=/root/waggle-heartbeat

[Install]
WantedBy=multi-user.target
