# arch/c1/extra/etc/systemd/system/waggle-firewall.service
[Unit]
Description=Setup Routing and Firewall

[Service]
Type=oneshot
ExecStart=/root/waggle-firewall

[Install]
WantedBy=multi-user.target
