# upstart/waggle_heartbeat.conf
start on (local-filesystems)
stop on runlevel[016]
respawn
chdir /usr/lib/waggle/waggle_image/
exec /usr/lib/waggle/waggle_image/heartbeat.sh

